# ev_smart_charger.service
[Unit]
Description=EV Smart Charger Optimizer
After=network.target

[Service]
Type=simple
User=peccz
WorkingDirectory=/home/peccz/ev_smart_charger
# Use the Python executable from the virtual environment
ExecStart=/home/peccz/ev_smart_charger/venv/bin/python src/main.py
Restart=on-failure
# Use standard syslog identifiers instead of direct file writing to avoid permission issues
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ev_smart_charger

[Install]
WantedBy=multi-user.target